# provision/roles/api/templates/api.conf.j2
[program:api]
command=python /home/pi/homeautomation/api/api.py
directory=/home/pi/homeautomation/api/
