# provision/roles/lights/templates/lights.conf.j2
[program:lights]
command=python /home/pi/homeautomation/lights/lights.py
directory=/home/pi/homeautomation/lights/
